# ngnix-flask-app/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn flask-app:app -b localhost:5000
directory=app